# service_files/xyz.openbmc_project.bios_bmc_smm_error_logger.service.in
[Unit]
Description=BIOS BMC SMM Error Logger

[Service]
Type=simple
ExecStart=@BINDIR@/bios-bmc-smm-error-logger
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
